# supervisord/conf.d/nginx.conf
[program:nginx]
command = /app/scripts/run_nginx.sh
autorestart = true
