# supervisord/conf.d/nginx.conf
[program:nginx]
command = /app/scripts/run-nginx.sh %(ENV_NGINX_ARGS)s
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
